# build/systemd/hive-ui.service
[Unit]
Description=Hive Web UI
PartOf=hive.target
After=hive-awsgw.service

[Service]
Type=simple
User=hive
Group=hive
ExecStart=/usr/local/bin/hive service hive-ui start
Restart=on-failure
RestartSec=5
OOMScoreAdjust=-500

[Install]
WantedBy=hive.target
