# build/systemd/hive-nats.service
[Unit]
Description=Hive NATS Message Bus
PartOf=hive.target
Before=hive-predastore.service hive-daemon.service hive-awsgw.service hive-vpcd.service

[Service]
Type=simple
User=hive
Group=hive
ExecStart=/usr/local/bin/hive service nats start
Restart=on-failure
RestartSec=5
OOMScoreAdjust=-500
Environment=HIVE_CONFIG_PATH=/etc/hive/nats/nats.conf

[Install]
WantedBy=hive.target
